# tang-server/usr/lib/systemd/system/tang-server.socket
[Unit]
Description=Tang Server socket
Documentation=man:tang(8)

[Socket]
ListenStream=9100
Accept=yes

[Install]
WantedBy=sockets.target
