# tang-server/usr/lib/systemd/system/tang-server-proxy.socket
[Unit]
Description=Tang Server proxy

[Socket]
ListenStream=9100

[Install]
WantedBy=sockets.target
