# tang-server/usr/lib/systemd/system/tang-server@.service.d/10-fix-container-name.conf
# systemd socket instance names are not valid Podman container
# names. Although seemingly undocumented, the first part of the
# instance is a sequence number, just strip off the rest and just use
# that. We don't need long term unique, so the risk is low.

[Service]
ExecStart=
ExecStart=/bin/sh -c '/usr/bin/podman run -i --name systemd-%p_$(echo "%i" | cut -d"-" -f1) --cidfile=%t/%N.cid --replace --rm --cgroups=split --network none -v systemd-tang-server:/var/db/tang:Z localhost/tang-server'
